# utils/v18.06.8/files/etc/uci-defaults/14_banner
#!/bin/sh

mv /etc/banner /etc/banner_tmp

# add AXIS Gear ascii art
AXISGEAR="
 █████╗ ██╗  ██╗██╗███████╗     ██████╗ ███████╗ █████╗ ██████╗
██╔══██╗╚██╗██╔╝██║██╔════╝    ██╔════╝ ██╔════╝██╔══██╗██╔══██╗
███████║ ╚███╔╝ ██║███████╗    ██║  ███╗█████╗  ███████║██████╔╝
██╔══██║ ██╔██╗ ██║╚════██║    ██║   ██║██╔══╝  ██╔══██║██╔══██╗
██║  ██║██╔╝ ██╗██║███████║    ╚██████╔╝███████╗██║  ██║██║  ██║
╚═╝  ╚═╝╚═╝  ╚═╝╚═╝╚══════╝     ╚═════╝ ╚══════╝╚═╝  ╚═╝╚═╝  ╚═╝ "

echo "$AXISGEAR" > /etc/banner

# old code to preserve bottom of old banner
#sed -e 9,14d -e 1,5d /etc/banner_tmp | cat >> /etc/banner

rm /etc/banner_tmp

